# roles/media_tools/files/compress-concerts.service
[Unit]
Description=Compress chaptered concert MKVs to HEVC for WiFi streaming
After=remux-concerts.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/compress_concerts.sh
Nice=19
IOSchedulingClass=idle
