# roles/media_tools/files/remux-concerts.service
[Unit]
Description=Remux BDMV concert rips into chaptered MKV for Jellyfin

[Service]
Type=oneshot
ExecStart=/usr/local/bin/remux_concerts.sh
Nice=19
IOSchedulingClass=idle
