# roles/media_tools/files/remux-concerts.timer
[Unit]
Description=Daily remux of new BDMV concert rips into chaptered MKV

[Timer]
OnCalendar=daily
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target
